# crontab/show_env_vars.sh
#!/usr/bin/env bash

/usr/bin/env > /tmp/env_vars.txt
